# supervisord.conf
[supervisord]
nodaemon=true

[program:python_pidisplayboard]
command=/usr/bin/python push_images.py
