# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:python_pidisplayboard]
command=cd /pidisplayboard; /bin/python push_images.py